# config/systemd/user/offlineimap-oneshot.service.d/override.conf
[Service]
Environment="NOTMUCH_CONFIG=/home/chris/.config/notmuch/config"
ExecStartPost=/usr/bin/notmuch new
