# config/systemd/user/offlineimap-oneshot.timer.d/override.conf
[Timer]
OnStartupSec=5min
OnUnitInactiveSec=15min
