# config/systemd/user/offlineimap-oneshot.timer.d/override.conf
[Timer]
OnBootSec=15min
OnUnitInactiveSec=30min
